# freshrss/rootfs/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant App: FreshRSS
# Runs the PHP-FPM daemon
# ==============================================================================

# Execute PHP-FPM
bashio::log.info "Starting php-fpm...."

exec php-fpm83 -R --nodaemonize
